# api/packer/api.service
[Unit]
Description=api

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/app/current/api

[Install]
WantedBy=multi-user.target
